# .config/X11/profile.d/01-xkb.sh
#!/usr/bin/env bash
## ek9/x11-config - https://github.com/ek9/x11-config
## 01-xkb.sh
## This file sets up xkb (keyboard keymap)

# Set xkbmap
#setxkbmap -model pc104 -layout us,nl -option grp:shifts_toggle &
setxkbmap -model pc104 -layout us &

