# .config/X11/profile.d/05-programs.sh
#!/usr/bin/env bash
## ek9/x11-config - https://github.com/ek9/x11-config
## 05-programs.sh
## This file sets up default programs for X11

[[ -x $(command -v firefox) ]] && export BROWSER=firefox

