# .config/X11/profile.d/80-gnome-keyring.sh
#!/usr/bin/env bash
## ek9/x11-config - https://github.com/ek9/x11-config
## 80-gnome-keyring.sh
## This file sets up gnome keyring

# Export session for gnome-keyring
export DESKTOP_SESION=gnome

