# .config/X11/profile.d/90-dpms.sh
#!/usr/bin/env bash
## ek9/x11-config - https://github.com/ek9/x11-config
## 90-dpms.sh
## This file sets up dpms
xset s off -dpms
xset m 0 0

