# .config/bash/profile.d/52-xft-force-gtk-qt.sh
# force Xft on old GTK and QT applications
export GDK_USE_XFT=1
export QT_XFT=true

