# .config/bash/profile.d/55-qt-theme.sh
# Make QT use GTK theme
export QT_STYLE_OVERRIDE=gtk
