# .config/bash/profile.d/60-wine-32bit-mode.sh
# use win32 for wine
export WINEARCH=win32

