# .config/bash/profile.d/70-dota2.sh
# Dota 2 font rendering fix
export MESA_GL_VERSION_OVERRIDE=2.1

